# parts/linux/cloud-init/artifacts/node-exporter/baseline/etc/systemd/system/node-exporter-restart.path
[Path]
PathModified=/var/lib/kubelet/pki/kubelet.crt

[Install]
WantedBy=multi-user.target
